# configs/ai-menshen.service
[Unit]
Description=AI-Menshen Service
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/tmp
ExecStart=/home/YOUR_USER/.local/bin/ai-menshen -config /home/YOUR_USER/.config/ai-menshen/config.toml
Restart=on-failure

[Install]
WantedBy=multi-user.target
